# installer/macos/com.sayzo.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sayzo.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Sayzo Agent.app/Contents/MacOS/sayzo-agent</string>
        <string>service</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Respawn on crashes only. A clean exit (exit code 0) from the tray
         "Quit" action must NOT trigger a relaunch — otherwise quitting from
         the menu bar does nothing visible, since launchd revives the
         process within seconds. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/sayzo-agent-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/sayzo-agent-stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
